# resources/btrfs-sync.timer.example
[Unit]
Description="Run btrfs-sync (wrapper) script hourly"

[Timer]
OnCalendar=*-*-* *:00:00
Unit=btrfs-sync.service

[Install]
WantedBy=multi-user.target
