# resources/btrfs-sync.service.example
[Unit]
Description="btrfs-sync-wrapper service unit file"

[Service]
ExecStart="/usr/local/bin/btrfs-sync-wrapper"
